# project/supervisor/conf/mongo-db.ini
[program:app]
command=/bin/bash -c "/usr/bin/supervisord -c /etc/supervisord.conf"
autorestart=false
